# spacepictl/etc/spacepi.target
[Unit]
Description=SpacePi Services Target

[Install]
WantedBy=multi-user.target
