# os/spacepi-ld-setup.service
[Unit]
Description=Library Path Initial Setup

[Service]
ExecStart=/spacepi/ld-setup.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
